# puppet/environment/vagrant/modules/mongodb/templates/mongod.conf.erb
#!upstart
description 'this file is enforced by puppet, and starts the mongod process.'

## start job defined in this file after system services, and processes have
##     already loaded (to prevent conflict).
start on runlevel [2345]

## stop upstart job
stop on runlevel [!2345]

## run upstart job as a background process
expect fork

## start upstart job
exec mongod --config /etc/mongod.conf